# psono-server.ini
[program:psono-server]
command = /usr/bin/gunicorn --bind 127.0.0.1:10100 psono.wsgi
directory=/home/psono/psono-server/psono
user = psono
autostart=true
autorestart=true
redirect_stderr=true
